# shell/profile.d/deno.sh
# Changing deno installation and package directory
export DENO_INSTALL="$HOME/.local/lib/deno"
export DENO_INSTALL_ROOT=$DENO_INSTALL

