# shell/profile.d/ghcup.sh
[ -f "/home/dan/.local/share/ghcup/env" ] && source "/home/dan/.local/share/ghcup/env" # ghcup-env

