# shell/profile.d/mocword.sh

export MOCWORD_DATA="$XDG_DATA_HOME/mocword/mocword.sqlite"

